# scripts/com.ga.expo-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ga.expo-tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/commando/Documents/G&amp;A User App/scripts/expo-tunnel.sh</string>
    <string>supervise</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/commando/Documents/G&amp;A User App</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/commando/Documents/G&amp;A User App/.expo/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/commando/Documents/G&amp;A User App/.expo/launchd.log</string>
</dict>
</plist>
